# deploy/fgpa-remote-lab.service
# /etc/systemd/system/fgpa-remote-lab.service on CT210
#
# Internal-only: nginx (see nginx-fgpa-remote-lab.conf) is the public
# entry point on port 8000 and reverse-proxies everything except
# hardware-lab traffic to this on 127.0.0.1:8001.
[Unit]
Description=FPGA Remote Lab backend (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fgpa-remote-lab/backend
ExecStart=/opt/fgpa-remote-lab/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
